# tbClientUpload.cron
*/10 * * * * /neousys/tbClient-igtInfo.sh
